# root/etc/services.d/baidupcs-web/run
#! /usr/bin/with-contenv bash

# 启动BaiduPCS-Go
exec s6-setuidgid emile  BaiduPCS-Go
